# .config/systemd/user/glava.service
[Unit]
Description=glava visualizer

[Service]
Type=simple
ExecStart=/usr/bin/glava --desktop
Restart=on-failure

[Install]
WantedBy=default.target
